# host/docker/buildroot/overlays/camera/etc/network/if-up.d/S41thttpd
#!/bin/sh

thttpd -p 8000 -d /root/www/ -l /var/log/thttpd.log
echo "HTTP server started"

