# host/docker/buildroot/overlays/camera/etc/network/if-up.d/S42start-camera
#!/bin/sh

echo " * START CAMERA ***"

/etc/network/if-up.d/gst-launch.sh $IFACE > /var/log/gst-launch 2>&1 &

